# bin/opensvc.init.debian
#!/bin/sh -e

### BEGIN INIT INFO
# Provides:             opensvc
# Required-Start:       $all
# Required-Stop:        
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    OpenSVC services startup script
### END INIT INFO

osvc_path=/opt/opensvc/bin
osvc_opts="--parallel"
allservices=${osvc_path}/allservices

case $1 in
start)
        ${osvc_path}/nodemgr pushasset
	${allservices} ${osvc_opts} boot
	;;
stop)
	${allservices} ${osvc_opts} shutdown
	;;
esac
