# bin/opensvc.init.debian
#!/bin/sh -e

### BEGIN INIT INFO
# Provides:             opensvc
# Required-Start:       $networking
# Required-Stop:        
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    OpenSVC services startup script
### END INIT INFO

osvc_path=/opt/opensvc/bin
allservices=${osvc_path}/allservices
allprimaryservices=${osvc_path}/allprimaryservices

case $1 in
start)
	${allprimaryservices} start
	${allservices} startstandby
	;;
stop)
	${allservices} stop
	;;
esac
