# bin/opensvc.init.redhat
#!/bin/bash
#
#	/etc/rc.d/init.d/opensvc
#
# Starts the services driven by OpenSVC
#
# chkconfig: 345 99 01
# description: Starts the services driven by OpenSVC whose
#              autostart node is this node.
# processname:

# Source function library.
. /etc/init.d/functions

osvc_path=/opt/opensvc/bin
osvc_opts="--parallel"
allservices=${osvc_path}/allservices

case $1 in
start)
	${allservices} ${osvc_opts} boot
	;;
stop)
	${allservices} ${osvc_opts} shutdown
	;;
esac
