# bin/init/darwin.com.opensvc.svcmgr.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.opensvc.svcmgr</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/share/opensvc/bin/init/opensvc.init.Darwin</string>
		<string>start</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/var/log/opensvc/svcmgr_boot_stderr.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/opensvc/svcmgr_boot_stdout.log</string>
</dict>
</plist>
