# bin/init/systemd.opensvc-agent.service
[Unit]
Description=OpenSVC Agent
Documentation=http://docs.opensvc.com/ file:/usr/share/doc/opensvc/
After=network.target network-online.target

[Service]
Type=idle
Environment="PATH=/opt/opensvc/bin:/opt/opensvc/etc:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/share/opensvc/bin/init/opensvc.init start
ExecStop=/usr/share/opensvc/bin/init/opensvc.init stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
