# bin/init/systemd.opensvc-agent.service
[Unit]
Description=OpenSVC agent
Documentation=https://docs.opensvc.com/ file:/usr/share/doc/opensvc/ man:nodemgr(1) man:svcmgr(1) man:svcmon(1)
Before=opensvc-services.service
After=network-online.target
After=time-sync.target
After=multi-user.target
After=network.target
After=docker.service
After=libvirtd.service
After=virt-guest-shutdown.target
After=libvirt-guests.service
After=blk-availability.service
Wants=libvirtd.service
Wants=libvirt-guests.service
Wants=blk-availability.service
Requires=opensvc-services.service

[Service]
Type=forking
TimeoutStopSec=5m
ExecStart=/usr/share/opensvc/bin/om daemon start
ExecStop=/usr/share/opensvc/bin/om daemon stop
KillMode=process
Restart=on-failure
OOMScoreAdjust=-1000
PIDFile=/var/lib/opensvc/osvcd.pid

[Install]
WantedBy=multi-user.target
