# bin/init/systemd.opensvc-agent.service
[Unit]
Description=OpenSVC Agent
Documentation=http://docs.opensvc.com/ file:/usr/share/doc/opensvc/
After=network.target network-online.target time-sync.target multi-user.target

[Service]
Type=forking
ExecStart=/usr/share/opensvc/bin/init/opensvc.init start
ExecStop=/usr/share/opensvc/bin/init/opensvc.init stop
KillMode=process
Restart=on-failure
OOMScoreAdjust=-1000

[Install]
WantedBy=multi-user.target
