# bin/init/systemd.opensvc-services.service
[Unit]
Description=OpenSVC agent and services oneshot stopper
Documentation=https://docs.opensvc.com/ file:/usr/share/doc/opensvc/ man:nodemgr(1) man:svcmgr(1) man:svcmon(1)
After=opensvc-agent.service

[Service]
Type=oneshot
TimeoutStopSec=1h
RemainAfterExit=true
ExecStart=/usr/share/opensvc/bin/om node pushasset
ExecStop=/usr/share/opensvc/bin/om daemon shutdown

